# task_solver/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

cd /app
export TASK_SOLVER_DB_PATH="/data/task_solver.db"
exec uvicorn app.main:app --host 0.0.0.0 --port 8099
